# wine-meta/etc/profile.d/wine-meta.sh
export PATH=/opt/wine-devel/bin:/opt/wine-staging/bin:/opt/wine-stable/bin:$PATH